# post_provision/etc/profile.d/path.sh
# Auto generated by /var/lib/apps/bin/post_provisioning
export PATH=$PATH:/var/lib/apps/bin
