# post_provision/etc/systemd/system/update-window.service
[Unit]
Description=Reboot if an update has been downloaded

[Service]
ExecStart=/bin/bash /opt/bin/update-window.sh
