# post_provision/etc/systemd/system/update-window.timer
[Unit]
Description=Reboot timer

[Timer]
OnCalendar=*-*-* 05,06:00/30:00
